# etc/cron.d/drivers-linux-firmware
@reboot root /usr/bin/drivers-linux-firmware
1 * * * * root /usr/bin/drivers-linux-firmware
